# Docker/root/etc/services.d/guacd/run
#!/bin/bash

echo "Starting guacamole guacd..."
s6-setuidgid root guacd -f